# configure-dotnet.service
[Unit]
Description=Configure bluetooth

[Service]
Type=oneshot
ExecStart=/opt/scripts/install_dotnet.sh

[Install]
WantedBy=multi-user.target